# systemd/paperless-backup.timer
[Unit]
Description=Paperless Backup – täglicher Cron-Ersatz

[Timer]
OnCalendar=*-*-* 03:00:00
Persistent=true
Unit=paperless-backup.service

[Install]
WantedBy=timers.target
